# products/music-streaming-agent/scripts/install/systemd/mpd.service.d/wait-for-network.conf
# MPD binds to a specific LAN IP in /etc/mpd_local.conf. If mpd.service starts
# before Wi-Fi has that address, bind() fails and the daemon exits.
# Ordering after network-online avoids that race (NetworkManager-wait-online on this Pi).
[Unit]
After=network-online.target
Wants=network-online.target
